# roles/local_backup_directory/templates/directory-backup.cron.j2
# Perform backup during the night
05 03 * * * root /usr/local/bin/{{ directory.split("/")[-1] }}-directory-backup.sh
