# roles/local_backup_mysql/templates/mysql-backup.cron.j2
# Performa backup during the night
05 03 * * * root /usr/local/bin/mysql-{{ database }}-backup.sh
